# etc/init.d/cryptmount
#!/sbin/openrc-run

description="Manages encrypted filesystems and swap with cryptmount"

depend() {
	need localmount
	use crypttab
	provide cryptmount
}

start_pre() {
	if [ ! -x /usr/bin/cryptmount ]; then
		eerror "/usr/bin/cryptmount not found"
		return 1
	fi
	if [ ! -f /etc/cryptmount/cmtab ]; then
		eerror "/etc/cryptmount/cmtab not found"
		return 1
	fi

	modprobe -q dm-crypt
}

start() {
	ebegin "Mounting cryptmount targets"
	/usr/bin/cryptmount --system-boot
	eend $?
}

stop() {
	ebegin "Unmounting cryptmount targets"
	/usr/bin/cryptmount --system-shutdown
	/usr/bin/cryptmount --safetynet
	eend $?
}
