# etc/init.d/dmcrypt
#!/sbin/openrc-run

description="Sets up encrypted devices from /etc/crypttab"

CRYPTTAB_FILE="/etc/crypttab"

depend() {
	need dmsetup
	before fsck localmount swap
	provide crypttab
}

start_pre() {
	[ -f "${CRYPTTAB_FILE}" ] || return 0
	modprobe -q dm-crypt
}

start() {
    [ ! -f "${CRYPTTAB_FILE}" ] && return 0
    ebegin "Setting up encrypted devices from ${CRYPTTAB_FILE}"

    local fail_marker="/tmp/.crypttab.fail.$$"
    rm -f "${fail_marker}"

    sed 's/#.*//;/^$/d' "${CRYPTTAB_FILE}" | \
    while IFS= read -r target source key_file options; do
        [ -z "${target}" ] || [ -z "${source}" ] && continue

        if cryptsetup status "${target}" >/dev/null 2>&1; then
            veinfo "  ${target} is already active"
            continue
        fi

        einfo "  Activating ${target}"

        [ "${key_file}" = "-" ] && key_file="none"

        if [ -n "${options}" ]; then
            # shellcheck disable=SC2086
            set -- $options
        else
            set --
        fi

        if ! cryptsetup open --key-file "${key_file}" --try-empty-password "$@" "${source}" "${target}"; then
            ewarn "    Failed to activate ${target}"
            touch "${fail_marker}"
        fi
    done

    if [ -f "${fail_marker}" ]; then
        rm -f "${fail_marker}"
        eend 1 "One or more encrypted devices failed to start"
        return 1
    else
        eend 0
        return 0
    fi
}

stop() {
	[ ! -f "${CRYPTTAB_FILE}" ] && return 0
	ebegin "Closing encrypted devices from ${CRYPTTAB_FILE}"

	tac "${CRYPTTAB_FILE}" | sed 's/#.*//;/^$/d' | \
	while IFS= read -r target _; do
		[ -z "${target}" ] && continue

		if cryptsetup status "${target}" >/dev/null 2>&1; then
			einfo "  Deactivating ${target}"
			cryptsetup close "${target}" 2>/dev/null || ewarn "    Failed to close ${target}"
		fi
	done

	eend 0
}
