# etc/init.d/dmeventd
#!/sbin/openrc-run

description="Starts the Device Mapper event daemon (dmeventd)"
command="/usr/sbin/dmeventd"
supervisor=supervise-daemon

depend() {
	after lvm dmsetup
}
