# etc/init.d/dmeventd
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

description="Starts the Device Mapper event daemon (dmeventd)"
command="/usr/sbin/dmeventd"
pidfile="/run/dmeventd.pid"

depend() {
	# Should start after lvm and device-mapper are ready.
	after lvm device-mapper
}
