# debian/postinst
#! /bin/bash

#   Add cryptsetup service files to runlevel/boot

/bin/cp /etc/init.d/{cryptdisks,cryptdisks-early} /etc/runlevels/boot || true

#   Update services list

rc-update -u
