# debian/postinst
#! /bin/bash

#   Update services list

rc-update -u
